# resize-fs.service
[Unit]
Description=Resize-FS one time run.

[Service]
Type=simple
ExecStart=/usr/bin/resize-fs

[Install]
WantedBy=default.target